# webhook-deploy.service
[Unit]
Description=GitHub Webhook Deployment Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
Environment="GITHUB_WEBHOOK_SECRET=your_secret_here"
ExecStart=/usr/bin/python3 /usr/local/bin/webhook-server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
